# docker/services/services.ini
; services.ini — process definitions for the dexbox container
; All services run as dexbox (non-root)

[program:xvfb]
command=/usr/bin/Xvfb :%(ENV_DISPLAY_NUM)s -ac -screen 0 %(ENV_WIDTH)sx%(ENV_HEIGHT)sx24 -retro -dpi 96 -nolisten tcp -nolisten unix
user=dexbox
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mutter]
command=/usr/bin/dbus-run-session -- bash -c "export XDG_RUNTIME_DIR=/tmp/runtime-user XDG_SESSION_TYPE=x11; mutter --display :%(ENV_DISPLAY_NUM)s --replace --sm-disable"
user=dexbox
environment=XDG_RUNTIME_DIR="/tmp/runtime-user",XDG_SESSION_TYPE="x11"
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:lxpanel]
command=/usr/bin/lxpanel --profile default
user=dexbox
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s"
priority=25
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tigervnc]
command=bash -c "/usr/bin/x0vncserver -fg -display :${DISPLAY_NUM} -rfbport 5900 -SecurityTypes=None -AlwaysShared=1 -localhost no --I-KNOW-THIS-IS-INSECURE -desktop dexbox-vnc -verbose -Log=*:stderr:100"
user=dexbox
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/tigervnc.log
stderr_logfile=/var/log/supervisor/tigervnc.err.log

[program:novnc]
command=bash -c "/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:6080 --web /opt/noVNC"
user=dexbox
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err.log

[program:ffmpeg-recorder]
command=ffmpeg -f x11grab -r 15 -s %(ENV_WIDTH)sx%(ENV_HEIGHT)s -i :%(ENV_DISPLAY_NUM)s -c:v libx264 -preset ultrafast -crf 32 /tmp/recordings/recording_%(ENV_DISPLAY_NUM)s.mp4
user=dexbox
autostart=false
autorestart=false
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s"

[program:rdp-keepalive]
command=/app/docker/bin/rdp-keepalive-wrapper.sh
directory=/app
user=dexbox
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s",HOME="/home/dexbox"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/rdp-keepalive.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/rdp-keepalive-error.log
stderr_logfile_maxbytes=10MB
priority=45
startsecs=3
startretries=999999999
stopwaitsecs=10


[program:dexbox-go-server]
command=/usr/local/bin/dexbox server
user=dexbox
autostart=true
autorestart=true
priority=59
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
